# crontab
*/30 * * * * /root/Desktop/workspace/GithubPush/Gitadd.py >/dev/null 2>&1
*/30 * * * * sleep 15; /root/Desktop/workspace/GithubPush/Gitrem.py >/dev/null 2>&1

